# deploy/launchd/com.rimrim.hunt2026-publish.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rimrim.hunt2026-publish</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>cd /Users/kristenho/projects/alpha-lab &amp;&amp; .venv/bin/python scripts/paper_publish.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/kristenho/projects/alpha-lab</string>
    <!-- hourly; the publisher commits ONLY when STATUS.md changed, so overnight (static equity)
         produces no commits. Widen to e.g. 7200 if you want fewer refreshes. -->
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/kristenho/projects/alpha-lab/artifacts/hunt2026/paper/publish.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kristenho/projects/alpha-lab/artifacts/hunt2026/paper/publish.log</string>
</dict>
</plist>
